# ansible/roles/web/templates/web.service.j2
[Unit]
Description=URL Shortener Web Application
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/app/app
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/app
Environment=UV_CACHE_DIR=/app/.uv_cache
ExecStartPre=/bin/sleep 2
ExecStart=/usr/local/bin/uv run uvicorn main:app --host 0.0.0.0 --port {{ app_port | default(8000) }}
Restart=always
RestartSec=3

# Environment variables
EnvironmentFile=/app/.env

# Security settings (temporarily relaxed for debugging)
# NoNewPrivileges=true
# ProtectSystem=strict
# ProtectHome=true
ReadWritePaths=/app

[Install]
WantedBy=multi-user.target